# deploy/iot-cx-agent.service
[Unit]
Description=IOT Cx Edge Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=swadmin
Group=swadmin
EnvironmentFile=-/etc/iot-cx-agent/edge-agent.env
WorkingDirectory=/home/swadmin/iot-cloud-commissioning
ExecStart=/home/swadmin/iot-cloud-commissioning/edge-agent/.venv/bin/iot-cx-agent --config /etc/iot-cx-agent/agent.yaml
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
